# supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=/go/bin/message_broker

[program:client]
command=/go/bin/client
